# ops/launchd/com.cockpit.dispatch-claude.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!-- Keeps the persistent `dispatch-claude` tmux session alive (interactive
       Claude Code, bypass permissions) — the claude-tmux adapter's target.
       Re-checked every 5 min; recreates the session if Claude Code exited. -->
  <key>Label</key><string>com.cockpit.dispatch-claude</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>exec /Users/agentsmyth/workspaces/dev/cockpit/ops/bin/ensure-dispatch-claude.sh</string>
  </array>
  <key>StartInterval</key><integer>300</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/agentsmyth/Library/Logs/cockpit-dispatch/claude-launchd.log</string>
  <key>StandardErrorPath</key><string>/Users/agentsmyth/Library/Logs/cockpit-dispatch/claude-launchd.log</string>
</dict>
</plist>
